# kitchen/site-cookbooks/rails_server/templates/default/profile.d/rails_app.sh.erb
export RAILS_ENV=<%= node[:app_config][:environment] %>